# deploy/writer-witness/writer-witness.service
[Unit]
Description=Trading Bot WebApp Writer Witness
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=writer-witness
Group=writer-witness
WorkingDirectory=/opt/trading-bot-witness/active/release
EnvironmentFile=/etc/trading-bot-witness/runtime.env
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONPYCACHEPREFIX=/dev/null
Environment=PATH=/usr/sbin:/usr/bin:/sbin:/bin
UnsetEnvironment=PYTHONPATH PYTHONHOME PYTHONSTARTUP PYTHONINSPECT PYTHONUSERBASE BASH_ENV ENV SHELLOPTS BASHOPTS CDPATH GLOBIGNORE LD_PRELOAD LD_LIBRARY_PATH LD_AUDIT LD_DEBUG LD_DEBUG_OUTPUT LD_PROFILE LD_SHOW_AUXV LD_BIND_NOW LD_BIND_NOT LD_ORIGIN_PATH LD_DYNAMIC_WEAK LD_HWCAP_MASK GLIBC_TUNABLES
ExecStart=/opt/trading-bot-witness/active/venv/bin/python -I -B -X utf8 -X pycache_prefix=/dev/null -m uvicorn writer_witness_app:app --host 127.0.0.1 --port 8011 --workers 1 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=3s
TimeoutStartSec=30s
TimeoutStopSec=15s
KillSignal=SIGTERM
UMask=0077

# The process is a private control-plane API. It needs only loopback network
# access to local PostgreSQL and the loopback Nginx reverse proxy.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadOnlyPaths=/opt/trading-bot-witness/active/release
ReadOnlyPaths=/etc/trading-bot-witness

[Install]
WantedBy=multi-user.target
