# deploy/writer-witness/writer-witness-activation-watchdog.service
[Unit]
Description=Reconcile an interrupted Writer Witness activation
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/writer-witness-activation-watchdog
UMask=0077
