# deploy/writer-witness/writer-witness-offsite-backup.service
[Unit]
Description=Encrypt and upload the latest Writer Witness backup off-host
After=network-online.target writer-witness-backup.service
Wants=network-online.target
ConditionPathExists=/etc/trading-bot-witness/offsite-backup.env

[Service]
Type=oneshot
EnvironmentFile=/etc/trading-bot-witness/offsite-backup.env
Environment=PATH=/usr/sbin:/usr/bin:/sbin:/bin
UnsetEnvironment=PYTHONPATH PYTHONHOME PYTHONSTARTUP PYTHONINSPECT PYTHONUSERBASE BASH_ENV ENV SHELLOPTS BASHOPTS CDPATH GLOBIGNORE LD_PRELOAD LD_LIBRARY_PATH LD_AUDIT LD_DEBUG LD_DEBUG_OUTPUT LD_PROFILE LD_SHOW_AUXV LD_BIND_NOW LD_BIND_NOT LD_ORIGIN_PATH LD_DYNAMIC_WEAK LD_HWCAP_MASK GLIBC_TUNABLES
ExecStart=/bin/bash /usr/local/sbin/writer-witness-offsite-backup
UMask=0077
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/backups/trading-bot-witness
ReadOnlyPaths=/etc/trading-bot-witness
