# deploy/writer-witness/writer-witness-activation-watchdog.timer
[Unit]
Description=Continuously reconcile Writer Witness activation journals

[Timer]
OnBootSec=15s
OnUnitActiveSec=15s
AccuracySec=1s
Persistent=true
Unit=writer-witness-activation-watchdog.service

[Install]
WantedBy=timers.target
