# deploy/writer-witness/writer-witness-activation-recovery.service
[Unit]
Description=Recover an interrupted Writer Witness activation
After=local-fs.target
Before=nginx.service writer-witness.service writer-witness-backup.service writer-witness-offsite-backup.service

[Service]
Type=oneshot
UnsetEnvironment=PYTHONPATH PYTHONHOME PYTHONSTARTUP PYTHONINSPECT PYTHONUSERBASE LD_PRELOAD LD_LIBRARY_PATH
ExecStart=/usr/bin/python3.12 -I -S -B -X utf8 -X pycache_prefix=/dev/null /usr/local/sbin/writer-witness-activation recover-boot
UMask=0077

[Install]
RequiredBy=nginx.service writer-witness.service writer-witness-backup.service writer-witness-offsite-backup.service
